# crontab
@hourly certbot-auto renew --post-hook "certbot-post-hook-lighttpd yimgo.fr"
