# crontab
@hourly certbot-auto renew --no-boostrap --post-hook "certbot-post-hook-lighttpd yimgo.fr"
